# mod029eu/system-bck/var/lib/dpkg/info/ifplugd.conffiles
/etc/apm/scripts.d/ifplugd
/etc/ifplugd/action.d/ifupdown
/etc/ifplugd/ifplugd.action
/etc/init.d/ifplugd
/etc/pm/sleep.d/74ifplugd
