# clear029eu/system-bck/var/lib/dpkg/info/dnsmasq.list
/.
/etc
/etc/default
/etc/default/dnsmasq
/etc/dnsmasq.conf
/etc/dnsmasq.d
/etc/dnsmasq.d/README
/etc/init.d
/etc/init.d/dnsmasq
/etc/insserv.conf.d
/etc/insserv.conf.d/dnsmasq
/etc/resolvconf
/etc/resolvconf/update.d
/etc/resolvconf/update.d/dnsmasq
/etc/runit
/etc/runit/runsvdir
/etc/runit/runsvdir/default
/etc/sv
/etc/sv/dnsmasq
/etc/sv/dnsmasq/.meta
/etc/sv/dnsmasq/.meta/installed
/etc/sv/dnsmasq/finish
/etc/sv/dnsmasq/log
/etc/sv/dnsmasq/log/run
/etc/sv/dnsmasq/run
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dnsmasq.service
/lib/systemd/system/dnsmasq@.service
/usr
/usr/lib
/usr/lib/resolvconf
/usr/lib/resolvconf/dpkg-event.d
/usr/lib/resolvconf/dpkg-event.d/dnsmasq
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/dnsmasq.conf
/usr/share
/usr/share/bug
/usr/share/bug/dnsmasq
/usr/share/bug/dnsmasq/presubj
/usr/share/dnsmasq
/usr/share/dnsmasq/installed-marker
/usr/share/doc
/usr/share/doc/dnsmasq
/usr/share/doc/dnsmasq/changelog.Debian.gz
/usr/share/doc/dnsmasq/copyright
/usr/share/runit
/usr/share/runit/meta
/usr/share/runit/meta/dnsmasq
/usr/share/runit/meta/dnsmasq/installed
/var
/var/log
/var/log/runit
/var/log/runit/dnsmasq
/etc/sv/dnsmasq/log/supervise
/etc/sv/dnsmasq/supervise
